# openstack-hao/install-openstack.conf
description "Install OpenStack"
start on runlevel [2345]
task
exec /root/install-openstack.sh >/root/install-openstack.log 2>&1
post-start exec rm -f /etc/init/install-openstack.conf
